# atarashi-hako.service.j2
[Unit]
Description=Atarashi Hako - Container update utility

[Service]
ExecStart=/bin/bash /usr/local/bin/atarashi-hako.sh

[Install]
WantedBy=multi-user.target
