# atarashi-hako.timer.j2
[Unit]
Description=Atarashi Hako - container update utility

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
